# stow/xdg/dot-config/autostart/dunst.desktop
[Desktop Entry]
Name=dunst
Comment=Dunst Notification System
Exec=/home/andrew/dotfiles/scripts/dunst.sh
Terminal=false
Type=Application

