# stow/xdg/dot-config/autostart/dunst.desktop
[Desktop Entry]
Name=dunst
Comment=Dunst Notification System
Exec=/usr/bin/dunst
Terminal=false
Type=Application

