# stow/xdg/dot-config/autostart/nm-applet.desktop
[Desktop Entry]
Name=nm-applet
Comment=NetworkManager for the system tray
Exec=/usr/bin/nm-applet
Terminal=false
Type=Application
