# stow/xdg/dot-config/autostart/picom.desktop
[Desktop Entry]
Name=picom
Comment=Picom Compositor
Exec=/usr/bin/picom --experimental-backends -b
Terminal=false
Type=Application
