# stow/xdg/dot-config/autostart/redshift.desktop
[Desktop Entry]
Name=redshift
Comment=Redshift
Exec=/usr/bin/redshift
Terminal=false
Type=Application
