# stow/xdg/dot-config/autostart/sxhkd.desktop
[Desktop Entry]
Name=sxhkd
Comment=Simple X hotkey daemon
Exec=bash -c sxhkd
Terminal=false
Type=Application
Hidden=true
