# stow/xdg/dot-config/autostart/unclutter.desktop
[Desktop Entry]
Name=Unclutter
Comment=Unclutter
Exec=unclutter
Terminal=false
Type=Application
