# stow/xdg/dot-config/autostart/xscreensaver.desktop
[Desktop Entry]
Name=xscreensaver
Comment=Screensaver for X11
Exec=xscreensaver -nosplash
Terminal=false
Type=Application
Hidden=true
