# stow/xdg/dot-config/autostart/xscreensaver.desktop
[Desktop Entry]
Name=XScreensaver
Comment=Screensaver for Xorg
Exec=/usr/bin/xscreensaver
Terminal=false
Type=Application
