# stow/xdg/dot-config/autostart/xscreensaver.desktop
[Desktop Entry]
Name=xscreensaver
Comment=Screensaver for X11
Exec=xscreensaver
Terminal=false
Type=Application
