# stow/xfce4/dot-config/autostart/xfce4-autostart.desktop
[Desktop Entry]
Exec=xfce4-autostart
Type=Application
Terminal=false
Name=Startup Script for XFCE4
