# stow/mopidy/dot-config/systemd/user/mopidy.service
[Unit]
Description=Mopidy
After=network-online.target

[Service]
ExecStart=mopidy

[Install]
WantedBy=default.target 
