# stow/neovim/dot-config/systemd/user/neovim-nightly.service
[Unit]
Description=Download Neovim Nightly

[Service]
Type=simple
ExecStart=%h/.local/bin/nvim-linux.sh
