# stow/xscreensaver/dot-config/autostart/xscreensaver.desktop
[Desktop Entry]
Version=1.0
Name=Xscreensaver
Exec=xscreensaver -no-splash
Terminal=false
Type=Application
