# stow/systemd/dot-config/systemd/user/duckdns.timer
[Unit]
Description=Update DuckDNS

[Timer]
OnUnitActiveSec=5m

[Install]
WantedBy=timers.target
