# stow/systemd/dot-config/systemd/user/duckdns.service
[Unit]
Description=Update DuckDNS

[Service]
Type=oneshot
ExecStart=%h/dotfiles/scripts/duck.sh
